# templates/postfix_exporter.service.j2
[Unit]
Description=PostfixExporter

[Service]
TimeoutStartSec=0
User=postfix_exporter
ExecStart={{ postfix_exporter_bin_path }} {{ postfix_exporter_options }}

[Install]
WantedBy=multi-user.target
